# assets/supervisor/nafder.ini
[program:nafder]
command=/usr/local/bin/nafder -t -c /dev/console /var/run/nafder
autostart=true
autorestart=true
stdout_logfile=/var/log/nafder.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/var/log/nafder.log
stderr_logfile_maxbytes=100MB
stopasgroup=true
priority=1
